# razbot_Software_RPi/PiBotBlueInidigo/pibot.service
[Unit]
Description=robot upstart pibot

[Service]
Type=oneshot
ExecStart=/bin/sh -c "pibot-start"

[Install]
WantedBy=multi-user.target
